# AUR/nvidia-optimus/bumblebeed-resume-git/bumblebeed-resume.install
post_install() {
    printf "NOTE: Enable systemd service via command:\n"
    printf "systemctl enable bumblebeed-resume"
}

post_remove() {
    printf "NOTE: Disable systemd service via command:\n"
    printf "systemctl disable bumblebeed-resume"
}
